# init/buildcurl.conf
description "buildcurl"
start on (started docker)
stop on runlevel [!2345]

respawn

respawn limit 100 1

env IMAGE=buildcurl/buildcurl:latest
env NAME=buildcurl
env PORT=80

pre-start script
    /usr/bin/docker pull $IMAGE || true
    /usr/bin/docker rm $NAME || true
end script

exec /usr/bin/docker run -p $PORT:80 --env-file /etc/buildcurl/.env -v /var/run/docker.sock:/var/run/docker.sock --name $NAME $IMAGE

post-stop script
    /usr/bin/docker stop -t 2 $NAME || true
    /usr/bin/docker rm $NAME || true
    sleep 2
end script
